# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize
priority=10
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=20
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:media-generation-worker]
command=/bin/sh -lc 'php artisan queue:work %(ENV_MEDIA_GENERATION_QUEUE_CONNECTION)s --queue=%(ENV_MEDIA_GENERATION_QUEUE_NAME)s,default --sleep=%(ENV_MEDIA_GENERATION_QUEUE_SLEEP_SECONDS)s --tries=%(ENV_MEDIA_GENERATION_QUEUE_TRIES)s --backoff=%(ENV_MEDIA_GENERATION_QUEUE_BACKOFF_SECONDS)s --timeout=%(ENV_MEDIA_GENERATION_QUEUE_TIMEOUT_SECONDS)s --max-jobs=%(ENV_MEDIA_GENERATION_QUEUE_MAX_JOBS)s --max-time=%(ENV_MEDIA_GENERATION_QUEUE_MAX_TIME_SECONDS)s --memory=%(ENV_MEDIA_GENERATION_QUEUE_MEMORY_MB)s'
process_name=%(program_name)s_%(process_num)02d
priority=30
autostart=true
autorestart=true
numprocs=%(ENV_MEDIA_GENERATION_QUEUE_CONCURRENCY)s
user=www-data
startsecs=0
stopwaitsecs=%(ENV_MEDIA_GENERATION_QUEUE_STOPWAIT_SECONDS)s
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0